# scripts/systemd/firecracker-vm@.service
[Unit]
Description=Atlas Firecracker VM %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The jailer re-creates its per-launch artifacts inside the jail on EVERY start
# and aborts with EEXIST if any already exist — the device nodes it mknod()s
# (/dev/net/tun, /dev/kvm, ...) and the firecracker.pid it writes. The jail root
# persists on disk across stop/start, so a restart hits last boot's leftovers
# ("Failed to create /dev/net/tun via mknod: File exists", "Failed to open file
# /firecracker.pid: File exists"). Remove the jailer-owned dev/, the stale pid
# file, and the stale API socket so each start gets a clean jail to re-populate.
# The persistent rootfs/kernel/config are untouched. `-` prefix on dev/: a first
# boot has none yet, and that is not an error.
ExecStartPre=-/bin/rm -rf /var/lib/atlas/virtual-machines/%i/jail/firecracker/%i/root/dev
ExecStartPre=/bin/rm -f /var/lib/atlas/virtual-machines/%i/jail/firecracker/%i/root/firecracker.pid
ExecStartPre=/bin/rm -f /var/lib/atlas/virtual-machines/%i/jail/firecracker/%i/root/run/firecracker.socket
# Networking (netns + veth + in-namespace tap) must exist before the jailer
# joins the namespace; ExecStartPre runs to completion before ExecStart.
ExecStartPre=/var/lib/atlas/bin/vm-network-up.sh %i
# Run Firecracker under the jailer via the per-VM launcher that provision-vm.sh
# generated. The launcher (not this unit) assembles the jailer argv, because the
# --cgroup cpu.max value carries an internal space that systemd's ExecStart
# word-splitting would shatter. The launcher de-privileges Firecracker to the
# per-VM uid/gid, chroots it into the per-VM jail, applies cgroup-v2 memory/cpu
# caps + fd/file rlimits, and joins the VM's own network namespace. It `exec`s
# the jailer, so the jailer is this unit's main process (see KillMode below).
ExecStart=/var/lib/atlas/virtual-machines/%i/jailer-launch.sh
ExecStopPost=/var/lib/atlas/bin/vm-network-down.sh %i
StandardOutput=append:/var/lib/atlas/virtual-machines/%i/log/firecracker.log
StandardError=append:/var/lib/atlas/virtual-machines/%i/log/firecracker.log
Restart=always
RestartSec=5s
# The jailer is the unit's main process and Firecracker is its child; KillMode
# mixed sends SIGTERM to the jailer and SIGKILL to the whole control group, so
# the jailed Firecracker dies with the unit rather than being orphaned.
KillMode=mixed

[Install]
WantedBy=multi-user.target
